# cis/roles/ubuntu24_cis/templates/etc/systemd/system/tmp.mount.j2
{{ file_managed_by_ansible }}
[Unit]
Description=Temporary Directory /tmp
ConditionPathIsSymbolicLink=!/tmp
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target
After=swap.target

[Mount]
What=tmpfs
Where=/tmp
Type=tmpfs

Options={{ ubtu24cis_tmp_partition_mount_options | unique | join(',') }}

[Install]
WantedBy=local-fs.target
